# src/etc/profile.d/grub-zfs-fix.sh
# see https://www.csparks.com/BootFedoraZFS/index.html
export ZPOOL_VDEV_NAME_PATH=YES

